# .etc/cron.d/lego
11 11 * * 0 root /var/www/.sh/lego/renew.sh > /dev/null 2>&1
