# .etc/cron.d/top
* * * * * root /var/www/.sh/top/report-high-load-average.sh > /dev/null 2>&1
